# fim/yara/vnsoc-yara.timer
[Unit]
Description=VN-SOC YARA daily scan timer

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
