# vuln-scan/vnsoc-trivy.timer
[Unit]
Description=VN-SOC Trivy daily vulnerability scan

[Timer]
OnCalendar=*-*-* 02:30:00
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target
